# helper/deploy/helper.service
# systemd unit for the yuntiandeng.com PAW helper inference service.
# Install:  sudo cp helper/deploy/helper.service /etc/systemd/system/yuntiandeng-helper.service
# Edit the paths below to match where the repo is checked out and the venv lives.
[Unit]
Description=yuntiandeng.com PAW helper inference service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# This helper is now a paw-helper CONTENT PACK; the framework provides the server.
# WorkingDirectory + PAW_HELPER_CONTENT point at the content pack (this dir).
WorkingDirectory=/opt/yuntiandeng-helper/repo/helper
Environment=PAW_HELPER_CONTENT=/opt/yuntiandeng-helper/repo/helper
# Offload inference to the PAW API: the in-process local model serializes behind a
# lock, so concurrent requests collapse (4 concurrent ~92s vs ~1.7s via the API).
Environment=PAW_HELPER_INFERENCE_BACKEND=remote_infer
# CORS allow-list for the data endpoints (/ask,/feedback,/health). The shared
# backend serves the embeddable widget to all of Yuntian's sites, so each site's
# origin is listed here. /widget.js itself is public (not CORS-gated).
Environment=HELPER_ALLOWED_ORIGINS=https://yuntiandeng.com,https://www.yuntiandeng.com,https://neural-os.com,https://www.neural-os.com,https://programasweights.com,https://www.programasweights.com
Environment=HELPER_FEEDBACK_LOG=/var/lib/yuntiandeng-helper/feedback.jsonl
Environment=HELPER_QUERY_LOG=/var/lib/yuntiandeng-helper/queries.jsonl
# Optional: higher compile limits / private programs.
# Environment=PAW_API_KEY=paw_sk_...
# PAW model cache (so it persists across restarts and isn't in the service user's home).
Environment=PAW_CACHE_DIR=/var/lib/yuntiandeng-helper/paw-cache
# `paw-helper` is the framework CLI (pip install -e /opt/yuntiandeng-helper/paw-helper).
ExecStart=/opt/yuntiandeng-helper/venv/bin/paw-helper serve --host 127.0.0.1 --port 8088
Restart=on-failure
RestartSec=3
# Run as an unprivileged user that can read the repo and write the cache/log dirs.
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
